# deploy/vantaiphucloc-backend.service
[Unit]
Description=Vantaiphucloc FastAPI Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=vantaiphucloc
Group=vantaiphucloc
WorkingDirectory=/opt/vantaiphucloc/backend
EnvironmentFile=/etc/vantaiphucloc/.env
Environment="PYTHONPATH=/opt/vantaiphucloc/backend"
ExecStart=/opt/vantaiphucloc/venv/bin/uvicorn app.main:app \
    --workers 2 \
    --host 127.0.0.1 \
    --port 8000 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vantaiphucloc-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
